# breaker-ios/layout/Library/LaunchDaemons/breaker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<!-- basic -->
	<key>Label</key>
	<string>breaker</string>
	<key>Program</key>
	<string>/usr/bin/BRDaemon</string>
	<key>RunAtLoad</key>
	<true/>

	<!-- keep alive options -->
	<key>KeepAlive</key>
	<dict>
		<!-- restart if exit code is not zero -->
		<key>SuccessfulExit</key>
		<false />

		<!-- crash, restart me -->
		<key>Crashed</key>
		<true />
	</dict>

	
	<!-- log -->
	<key>StandardErrorPath</key>
	<string>/tmp/breaker_err.log</string>
	<key>StandardOutPath</key>
	<string>/tmp/breaker_std.log</string>


	<!-- debug -->
	<key>Debug</key>
    <true/>
    <key>SoftResourceLimits</key>
    <dict>
        <key>Core</key>
        <integer>9223372036854775807</integer>
    </dict>
    <key>HardResourceLimits</key>
    <dict>
        <key>Core</key>
        <integer>9223372036854775807</integer>
    </dict>


</dict>
</plist>

