# breaker-ios/layout/Library/LaunchDaemons/brwatchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<!-- basic -->
	<key>Label</key>
	<string>brwatchdog</string>
	<key>Program</key>
	<string>/usr/bin/BRWatchDog</string>
	<key>RunAtLoad</key>
	<true/>

	<!-- keep alive options -->
	<key>StartInterval</key>
	<integer>60</integer>
	
	<!-- log -->
	<key>StandardErrorPath</key>
	<string>/tmp/brwatchdog_err.log</string>
	<key>StandardOutPath</key>
	<string>/tmp/brwatchdog_std.log</string>


</dict>
</plist>

